# meta-iot2000-example/recipes-galileo/galileo-target/files/galileo-sketch-reset.service
[Unit]
Description=Galileo Arduino Sketch Reset

[Service]
Type=simple
ExecStart=/opt/cln/galileo/galileo_sketch_reset

[Install]
WantedBy=multi-user.target
