# meta-iot2000-example/recipes-galileo/galileo-target/files/galileo-target.service
[Unit]
Description=Galileo Arduino Layer

[Service]
Type=simple
ExecStart=/bin/sh /opt/cln/galileo/launcher.sh

[Install]
WantedBy=multi-user.target
